# app.conf
[program:reducer]
command=cargo watch -x run
directory=/app
autostart=true
autorestart=true
environment=PATH="/root/.cargo/bin:%(ENV_PATH)s"

